# daemon/com.openclaw.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.scheduler</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>system/scheduler/scheduler.py</string>
        <string>run</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ryleebenson/Desktop/OPENCLAW</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/ryleebenson/Desktop/OPENCLAW/logs/scheduler-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ryleebenson/Desktop/OPENCLAW/logs/scheduler-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/ryleebenson/Desktop/OPENCLAW</string>
    </dict>
</dict>
</plist>
